# ansible/roles/backend/templates/ai-backend.service.j2
[Unit]
Description=AI Backend (Node.js + Ollama proxy)
After=network.target ollama.service docker.service
Wants=ollama.service docker.service

[Service]
Type=simple
User={{ app_user }}
WorkingDirectory={{ app_dir }}

ExecStart=/usr/bin/node {{ app_dir }}/server.js

Restart=always
RestartSec=3

Environment=NODE_ENV=production
Environment=OLLAMA_MODEL={{ ollama_model }}

[Install]
WantedBy=multi-user.target
